# telewall/etc/crontab
@weekly /telewall/bin/import_blacklist.sh
@daily /telewall/bin/cleanup.sh
